# etc/s6-overlay/s6-rc.d/gravity-sync/up
/etc/s6-overlay/scripts/gravity-sync-start
